# conf/supervisor-ansible.ini
[program:ansible-container]
stopwaitsecs=30
directory=/home/ec2-user/app
command=/usr/local/bin/docker-compose up -d
